# templates/freshdesk-openstack-agent.service.erb
[Unit]
Description=freshdesk-openstack-agent
After=ntp.service

[Service]
<% if @rfc1918_gateway and @http_proxy -%>
Environment="HTTP_PROXY=<%= @http_proxy %>" "HTTPS_PROXY=<%= @http_proxy %>" "NO_PROXY=localhost,127.0.0.1,172.16.0.0/12,115.146.80/20,.test.rc.nectar.org.au,.dev.rc.nectar.org.au,.rc.nectar.org.au"
<% end -%>
User=<%= @user %>
Group=<%= @group %>
Type=simple
WorkingDirectory=<%= @base_dir %>
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /var/lock/nectar-freshdesk /var/log/nectar-freshdesk /var/lib/nectar-freshdesk
ExecStartPre=/bin/chown <%= @user %>:<%= @group %> /var/lock/nectar-freshdesk /var/log/nectar-freshdesk /var/lib/nectar-freshdesk
ExecStart=<%= @venv_dir %>/bin/freshdesk-openstack-agent
Restart=on-failure
LimitNOFILE=65535
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
