# files/upstart/freshdesk-pip.conf
#!upstart
#
# Freshdesk-PIP Upstart configuration
#
#
description "freshdesk-pip upstart configuration."
 
start on runlevel [2345]
stop on runlevel [06]


pre-start script
    # The minUptime and spinSleepTime settings stop Forever from thrashing if
    # the application fails immediately on launch. This is generally necessary
    # to avoid loading development servers to the point of failure every time
    # someone makes an error in application initialization code, or bringing
    # down production servers the same way if a database or other critical
    # service suddenly becomes inaccessible.
    mkdir -p /var/log/freshdesk /var/run/freshdesk
    chown -R freshdesk:freshdesk /var/log/freshdesk /var/run/freshdesk
    exec su -c ". ~/.nvm/nvm.sh && NODE_ENV=production forever --uid freshdesk-pip --pidFile /var/run/freshdesk/pip.pid -a  -l /var/log/freshdesk/pip.log --minUptime 5000 --spinSleepTime 2000 start /opt/freshdesk-pip/app.js" freshdesk
end script
 
post-stop script
    # Here we're using the post-stop script to stop the Node.js application
    # process so that Forever is given a chance to do its thing and tidy up
    # its data. Note that doing it this way means that each application that
    # runs under Forever must have a different start file name, regardless of
    # which directory it is in.
    exec su -c ". ~/.nvm/nvm.sh && forever stop freshdesk-pip" freshdesk
end script